# applications/luci-app-kge-vpn-watchdog/root/etc/uci-defaults/80_vpn_watchdog
#!/usr/bin/env ucode
'use strict';
// UCI defaults for vpn_watchdog (zero shell; ucode only).
// Creates /etc/config/vpn_watchdog and sets defaults. Writable site list under /etc (overlay).

let cursor = require('uci').cursor;
let fs = require('fs');

// Ensure config file exists (touch)
try {
	let f = fs.open('/etc/config/vpn_watchdog', 'a');
	if (f) f.close();
} catch (e) {}

let u = cursor();
u.load('vpn_watchdog');
let s = u.sections('vpn_watchdog', 'watchdog');
let sec = (s && s.length > 0) ? s[0]['.name'] : null;
if (!sec) {
	sec = u.add('vpn_watchdog', 'watchdog');
	if (!sec) { u.unload(); exit(1); }
}
u.set('vpn_watchdog', sec, 'enabled', '1');
u.set('vpn_watchdog', sec, 'log_path', '');
u.set('vpn_watchdog', sec, 'sites_file', '/etc/kge-vpn-watchdog/sites.conf');
u.set('vpn_watchdog', sec, 'vpn_iface', '');
u.set('vpn_watchdog', sec, 'polling_iface', '');
u.set('vpn_watchdog', sec, 'sleep_after_switch', '3');
u.set('vpn_watchdog', sec, 'curl_connect_timeout', '3');
u.set('vpn_watchdog', sec, 'curl_max_time', '8');
u.set('vpn_watchdog', sec, 'cron_enabled', '1');
u.set('vpn_watchdog', sec, 'cron_interval', '5');
u.set('vpn_watchdog', sec, 'cron_schedule', '');
u.set('vpn_watchdog', sec, 'run_timeout', '120');
u.set('vpn_watchdog', sec, 'log_tail_lines', '200');
u.set('vpn_watchdog', sec, 'peer_whitelist', '');
u.commit('vpn_watchdog');
u.unload();

// Create dir for writable sites (like mkdir -p)
try {
	let fp = fs.popen('mkdir -p /etc/kge-vpn-watchdog 2>/dev/null', 'r');
	if (fp) fp.close();
} catch (e) {}
exit(0);
